# files/a3s5a:and_his_name_was_bram_moolenaar.service
[Unit]                      
Description=clue2 

[Service] 
Type=oneshot 
ExecStart=/usr/bin/clue2 

[Install] 
WantedBy=multi-user.target
